# supervisord.conf
[supervisord]
nodaemon=true

[program:api_and_angularjs]
command=/bin/bash -l -c 'cd /app; /usr/local/bundle/bin/bundle exec puma -C config/puma.rb'
autorestart=true
